# etc/init.d/qdisc
#!/sbin/openrc-run
# (CAKE) Use flows for multiple network connections or flowblind for a single network connection.

description="Qdisc OpenRC Service"
command="/sbin/tc"
command_args="qdisc replace dev wlan0 root cake bandwidth 2Mbit besteffort flowblind nonat nowash no-ack-filter split-gso rtt 5ms raw overhead 0"
